# my.cron
*/1 * * * * /root/test/auto_failover.sh
*/3 * * * * /root/test/auto_failover.sh
*/5 * * * * /root/test/codelines.sh
